# systemd/huedinsteon.service
# SPDX-License-Identifier: Apache-2.0

[Unit]
Description=HuedInsteon bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/huedinsteon
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/huedinsteon/.venv/bin/huedinsteon run /etc/huedinsteon/config.toml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
